# v2/api.thewatcher.io/api.thewatcher.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn/socket
ListenStream=0.0.0.0:8013

[Install]
WantedBy=sockets.target
